# airemote-master/sdcard/200IDK-A2/profile.d/airemote.sh
export PATH=/usr/local/python3.7/bin:$PATH
export PYTHONPATH=/usr/local/python3.7/lib/python3.7/site-packages:$PYTHONPATH

. /usr/local/Ascend/ascend-toolkit/set_env.sh
. /etc/profile.d/color.sh

alias py=python3.7
alias python3=python3.7

